# pkg/rfTools/inst/crontabs/crontab_lin
# m h  dom mon dow   command
15,50 * * * * /home/rforge/lib/R/rfTools/exec/rf_build_packages >> /home/rforge/rf_build_packages.log 2>&1
45 11 * * * /home/rforge/lib/R/rfTools/exec/rf_update_R >> /home/rforge/rf_update_R.log 2>&1
